# overlays/home/ovos/.config/systemd/user/g4f.service
[Unit]
Description=GPT4Free

[Service]
ExecStart=%h/.venvs/ovos/bin/python -m g4f.api.run

[Install]
WantedBy=default.target
